# redmine.crontab
*/5 * * * * . /var/local/redmine/scripts/receive_imap_env.sh; /var/local/redmine/scripts/receive_imap.sh
0 6 * * * . /var/local/redmine/scripts/update-repositories.env.sh; /var/local/redmine/scripts/update-repositories.sh
